# kvm-media-player/kvm-connected.service
[Unit]
Description=KVM Connected

[Service]
Type=exec
RemainAfterExit=true
ExecStart=/usr/local/bin/on-kvm-connected
